# bot/src/main/jib/var/spool/cron/crontabs/root
# m h  dom mon dow   command
45 05 * * * /bin/sh /usr/local/bin/restartApp.sh
